# scripts/liquidity-bot/com.clawnet.liquidity-bot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawnet.liquidity-bot</string>

  <!-- Update this path to your repo root -->
  <key>WorkingDirectory</key>
  <string>/Users/yourname/workspace/private-repo/Bots/clawnet</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LIQUIDITY_BOT_CONFIG</key>
    <string>/Users/yourname/workspace/private-repo/Bots/clawnet/scripts/liquidity-bot/config.local.json</string>
    <key>LIQUIDITY_BOT_LOG_DIR</key>
    <string>/Users/yourname/workspace/private-repo/Bots/clawnet/scripts/liquidity-bot/logs</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yourname/workspace/private-repo/Bots/clawnet/scripts/liquidity-bot/run-scheduled.sh</string>
    <string>--count</string>
    <string>1</string>
    <string>--cycles</string>
    <string>2</string>
  </array>

  <!-- Every hour at minute 10 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/yourname/workspace/private-repo/Bots/clawnet/scripts/liquidity-bot/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yourname/workspace/private-repo/Bots/clawnet/scripts/liquidity-bot/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
